# tools/dockerize/webportal/usr/lib/systemd/system/daoliagent.service
[Unit]
Description=DaoliAgent Server
After=syslog.target network.target

[Service]
Restart=always
User=root
ExecStart=/usr/bin/daoliagent

[Install]
WantedBy=multi-user.target
